# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:hubot]
command=/src/bin/hubot --adapter hipchat